# src/test/data/deb/etc/init/hiapp.conf
description     "Hi server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

console none

pre-start script
    test -f /opt/hiapp/hiapp.jar || { stop; exit 0; }
    test -f /etc/hiapp.conf || { stop; exit 0; }
end script

script
    . /etc/hiapp.conf
    su hiapp -c "java -jar /opt/hiapp/hiapp.jar"
end script
